# chatroom/chat.bak
#!/bin/bash
#/etc/init.d/chat

### BEGIN INIT INFO
# Provides:          chat
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: chat
# Description:       chat pop & imap daemon
### END INIT INFO

case "$1" in
    start)
        echo "Starting chat service..."
        /home/yebin/pyprojects/chatroom/.venv/bin/python /home/yebin/pyprojects/chatroom/chat.py > /dev/null 2>&1 &
        ;;
    stop)
        echo "Stopping chat service..."
        killall /home/yebin/pyprojects/chatroom/.venv/bin/python
        ;;
    *)
        echo "Usage: /etc/init.d/chat{start|stop}"
        exit 1
        ;;
esac

exit 0
